# ser/battery-warning.service
[Unit]
Description=Battery Low Notification

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-warning.sh
